# roles/influxdb/templates/influxdb.service.j2
# This file is cotrolled by an ansible playbook, please do not edit!
# https://github.com/scline/ansible-home

[Unit]
Description=Influxdb
After=network.target docker.service
Requires=docker.service

[Install]
WantedBy=multi-user.target

[Service]
ExecStartPre=/usr/bin/docker pull influxdb:latest
# the - essentially ignores failure, it isn't made clear in the systemd manpage that
# this is the expected behavior, but it is. If this `-` prefixed command fails
# no other `ExecStartPre` commands will be executed, so we run the rm last.
ExecStartPre=-/usr/bin/docker rm -f influxdb
ExecStart=/usr/bin/docker run \
                          --name influxdb \
                          --init \
                          --pid=host \
                          --net=host \
                          -v /var/docker/influxdb:/var/lib/influxdb \
                          -e TZ={{ server_timezone }} \
                          influxdb:latest

ExecStop=/usr/bin/docker stop influxdb
# The rm of the container should be covered by the `--rm` flag
# in the run statement. However we want to make sure it's extra
# dead.
ExecStopPost=/usr/bin/docker rm -f influxdb
Restart=always
RestartSec=10
