# old-things/roles/netdata/files/netdata.service
# /etc/systemd/system/netdata.service
[Unit]
Description=Netdata Service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
StandardError=null
StandardOutput=null
WorkingDirectory=/etc/docker/compose/netdata
ExecStartPre=/usr/bin/docker-compose -f /etc/docker/compose/netdata/docker-compose.yml pull
ExecStart=/usr/bin/docker-compose -f /etc/docker/compose/netdata/docker-compose.yml up -d --force-recreate
ExecStop=/usr/bin/docker-compose -f /etc/docker/compose/netdata/docker-compose.yml down

[Install]
WantedBy=multi-user.target
