# roles/home-assistant/templates/home-assistant.service.j2
# This file is cotrolled by an ansible playbook, please do not edit!
# https://github.com/scline/ansible-home

[Unit]
Description=HomeAssistant
After=network.target docker.service
Requires=docker.service

[Install]
WantedBy=multi-user.target

[Service]
ExecStartPre=/usr/bin/docker pull homeassistant/home-assistant:latest
# the - essentially ignores failure, it isn't made clear in the systemd manpage that
# this is the expected behavior, but it is. If this `-` prefixed command fails
# no other `ExecStartPre` commands will be executed, so we run the rm last.
ExecStartPre=-/usr/bin/docker rm -f home-assistant
ExecStart=/usr/bin/docker run \
                          --name home-assistant \
                          --init \
                          --pid=host \
                          --device=/dev/ttyACM0:/zwaveusbstick \
                          -p 80:8123 \
                          -v /var/docker/home-assistant:/config  \
                          -e TZ=America/Los_Angeles \
                          homeassistant/raspberrypi3-homeassistant:latest

ExecStop=/usr/bin/docker stop home-assistant
# The rm of the container should be covered by the `--rm` flag
# in the run statement. However we want to make sure it's extra
# dead.
ExecStopPost=/usr/bin/docker rm -f home-assistant
Restart=always
RestartSec=10
